# root/etc/s6-overlay/s6-rc.d/init-mod-universal-flatpak/run
#!/usr/bin/with-contenv bash

echo "\
    flatpak" >> /mod-repo-packages-to-install.list

#flatpak remote-add --user --if-not-exists flathub https://flathub.org/repo/flathub.flatpakrepo
